# server/roles/loopback/templates/loopback.conf
start on started networking
stop on stopping networking

#respawn
#expect daemon

script
  touch /tmp/upstart
end script
